# launchd/com.disk-magician.gemini-dedup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.disk-magician.gemini-dedup</string>
  <key>ProgramArguments</key>
  <array>
    <string>@BASH@</string>
    <string>@REPO_ROOT@/scripts/symlink-shared-gemini.sh</string>
    <string>--clean</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AO_GEMINI_DEDUP_APPROVED</key>
    <string>1</string>
  </dict>
  <!--
    StartInterval, not StartCalendarInterval: see com.disk-magician.colima-prune.plist
    for the full rationale — the Sunday-04:30 calendar schedule never fired
    in 16+ days loaded because launchd drops a missed calendar window when
    asleep instead of queuing a catch-up run. 604800s = 7 days.
  -->
  <key>StartInterval</key>
  <integer>604800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/disk-magician-gemini-dedup.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/disk-magician-gemini-dedup.log</string>
</dict>
</plist>
