# launchd/com.disk-magician.hermes-vacuum.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly sweeper: reclaim WAL bloat from ~/.hermes/state.db.

  SAFETY (incident 2026-07-22, jleechan-4dtg): this job previously ran an
  unconditional `sqlite3 state.db 'VACUUM;'`, claiming unqualified safety
  "while daemon is up." On a 5.9 GB WAL-mode db, VACUUM writes a full
  second copy of the file -- it drove host free space from 12 GB to 5 GB
  in ~13 minutes and nearly wedged the disk, exactly when disk pressure
  is already the reason you'd want a cleanup sweeper to run. This plist
  now only passes --apply (no --full-vacuum), so the script runs
  `PRAGMA wal_checkpoint(TRUNCATE);` -- flushes the WAL into the main
  file and truncates state.db-wal, no second full-size copy, safe with
  the daemon up. Full VACUUM is opt-in only (scripts/vacuum_hermes_state.sh
  --apply --full-vacuum), free-space-guarded (refuses unless free space
  >= max(2x current DB size, 15 GB)), and intentionally NOT wired into
  this automated schedule -- run it manually when you want the file to
  actually shrink.

  Schedule:  Weekly (StartInterval, see below) — same family as the
             other state.db / filesystem hygiene jobs. StartInterval,
             not StartCalendarInterval: this job was never actually
             installed since it was added 2026-07-06 (the installer glob
             would have picked it up, but was never re-run). Its
             sibling early-AM calendar sweepers (colima-prune,
             gemini-dedup, playwright-dedup) also never fired once
             loaded — launchd drops a missed StartCalendarInterval
             window entirely if the Mac is asleep at that wall-clock
             moment, unlike cron. Switching to StartInterval avoids
             that failure mode. 604800s = 7 days.

  state.db is currently ~5.9 GB. Checkpointing prevents the WAL itself
  from growing large between writes; it does not shrink the main file
  (that needs VACUUM, run manually and guarded as above).

  Template uses @HOME@ placeholders so install scripts can substitute
  $HOME before copying into ~/Library/LaunchAgents/. Do not commit a
  fully-resolved path; keep @HOME@ literal in the repo.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.disk-magician.hermes-vacuum</string>
  <key>ProgramArguments</key>
  <array>
    <!--
      Use bash 5.x explicitly. @BASH@ on macOS is 3.2.57, which
      lacks the regex / array features the script depends on.
      launchd's default environment PATH is /usr/bin:/bin:/usr/sbin:/sbin,
      so "env bash" would also resolve to the wrong binary.
    -->
    <string>@BASH@</string>
    <string>@REPO_ROOT@/scripts/vacuum_hermes_state.sh</string>
    <string>--apply</string>
  </array>
  <!--
    StartInterval, not StartCalendarInterval: see comment block above
    for the full rationale. 604800s = 7 days.
  -->
  <key>StartInterval</key>
  <integer>604800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/disk-magician-hermes-vacuum.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/disk-magician-hermes-vacuum.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!--
      sqlite3 lives at /opt/homebrew/bin on Apple Silicon; macOS does
      not include it by default. Without this PATH, launchd would
      fail with "sqlite3: command not found" and exit 127.
    -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>