# launchd/com.disk-magician.colima-prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.disk-magician.colima-prune</string>
  <key>ProgramArguments</key>
  <array>
    <string>@BASH@</string>
    <string>@REPO_ROOT@/scripts/cleanup_colima.sh</string>
    <string>--clean</string>
    <string>--prune-volumes</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DOCKER_VOLUMES_APPROVED</key>
    <string>1</string>
  </dict>
  <!--
    StartInterval, not StartCalendarInterval: the prior Sunday-03:45
    calendar schedule never fired in 16+ days of being loaded (`launchctl
    print` showed "last exit code = (never exited)") because launchd drops
    a missed StartCalendarInterval window entirely if the Mac is asleep at
    that exact wall-clock moment — unlike cron, it does not queue a
    catch-up run. StartInterval instead fires as soon as possible after
    the elapsed-seconds window closes, including right after wake, which
    matches the observed-working pattern of the other repo-root sweepers
    (pressure-sweep, drilldown, downloads-evidence). 604800s = 7 days.
  -->
  <key>StartInterval</key>
  <integer>604800</integer>
  <key>StandardOutPath</key>
  <string>/tmp/disk-magician-colima-prune.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/disk-magician-colima-prune.log</string>
</dict>
</plist>
