# src/disk_magician/launchd/com.disk-magician.worktree-venvs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly sweeper: strip Python venvs from dormant Git worktrees.

  Schedule:  Sundays at 04:00 local time (Weekday=0, Hour=4, Minute=0).
  Action:    runs scripts/cleanup_worktree_venvs.sh --clean --min-age 14
             under WORKTREE_APPROVED=1, which satisfies the script's
             worktree-safety gate.

  Why WORKTREE_APPROVED is baked in:
  The script refuses --clean without WORKTREE_APPROVED=1 in the environment
  (matching the worktree-safety rule in user_scope CLAUDE.md, which flags
  any delete inside a worktree as requiring explicit approval). This
  scheduled job is itself that explicit approval: the operator chose to
  install this plist, and the plist runs the script with the approval
  token set, so the script's gate passes by design.

  Template uses @HOME@ placeholders so install scripts can substitute
  $HOME before copying into ~/Library/LaunchAgents/. Do not commit a
  fully-resolved path; keep @HOME@ literal in the repo.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.disk-magician.worktree-venvs</string>
  <key>ProgramArguments</key>
  <array>
    <!--
      Use bash 5.x explicitly. @BASH@ on macOS is 3.2.57, which lacks
      associative arrays (declare -A) that the script depends on.
      launchd's default environment PATH is /usr/bin:/bin:/usr/sbin:/sbin,
      so "env bash" would also resolve to the wrong binary.
    -->
    <string>@BASH@</string>
    <string>@REPO_ROOT@/scripts/cleanup_worktree_venvs.sh</string>
    <string>--clean</string>
    <string>--min-age</string>
    <string>14</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WORKTREE_APPROVED</key>
    <string>1</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/disk-magician-worktree-venvs.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/disk-magician-worktree-venvs.log</string>
</dict>
</plist>
